# deploy/roles/rpi-graphite/templates/graphite-api.socket.j2
# This file was generated by Ansible for {{ ansible_fqdn }}
# Do NOT modify this file by hand!
[Unit]
Description=graphite-api socket

[Socket]
ListenStream=/run/graphite-api.sock
ListenStream=127.0.0.1:{{ graphite_api_process_bind }}

[Install]
WantedBy=sockets.target
