# deploy/roles/rpi-graphite/templates/graphite-api.service.j2
# This file was generated by Ansible for {{ ansible_fqdn }}
# Do NOT modify this file by hand!

[Unit]
Description=Graphite-API service
Requires=graphite-api.socket

[Service]
ExecStart={{ gunicorn_path }} -w2 graphite_api.app:app
Restart=on-failure
#User=graphite
#Group=graphite
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
PrivateTmp=true

[Install]
WantedBy=multi-user.target
