# deploy/roles/rpi-graphite/templates/supervisor-graphite-api.conf.j2
# This file was generated by Ansible for {{ ansible_fqdn }}
# Do NOT modify this file by hand!
[program:graphite-api]
command = gunicorn -w2 graphite_api.app:app -b 127.0.0.1:{{ graphite_api_process_bind }}
autostart = true
autorestart = true
